# templates/upstart.conf
description "uWSGI started with upstart"

start on vagrant-mounted
stop on runlevel [016]

console log

exec {{home}}/bin/uwsgi --ini {{conf}}/uwsgi.ini
